# docker/rootfs/etc/systemd/system/user@.service.d/resources.conf
# https://lilac.readthedocs.io/en/latest/setup.html#configure-other-parts
[Service]
CPUWeight=100
OOMScoreAdjust=0
